# config/systemd/rackman.service
[Unit]
Description=RackMan Node.js Service Manager

[Service]
TimeoutStartSec=0
ExecStart=/usr/local/bin/node /usr/local/bin/rackman /web/MyProject